# scripts/terminaleyes-pi.service
# terminaleyes-pi.service — Systemd unit for the Pi keyboard REST API.
#
# This is a TEMPLATE. The deploy_pi.sh script generates the actual
# service file with the correct paths for the target user.
#
# Manual install:
#   sudo cp scripts/terminaleyes-pi.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable terminaleyes-pi

[Unit]
Description=TerminalEyes Pi Keyboard REST API
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/andras/terminaleyes
ExecStartPre=/bin/bash /home/andras/terminaleyes/scripts/setup_usb_gadget.sh ecm
ExecStartPre=/bin/bash /home/andras/terminaleyes/scripts/radio_mode.sh apply
ExecStart=/home/andras/terminaleyes/.venv/bin/terminaleyes-pi
Restart=on-failure
RestartSec=5
Environment=PATH=/home/andras/terminaleyes/.venv/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
